# ansible-infra/roles/galera_check/templates/galerachk.socket.j2
[Unit]
Description=Galera Cluster node check socket

[Socket]
ListenStream={{ galera_check_port }}
Accept=yes

[Install]
WantedBy=sockets.target
